# runner/runner.service
# Brings the self-hosted runner up at boot.
#
# The runner is a docker compose service someone has to `up`. Without this unit, every reboot leaves the
# platform with no runner — and a deploy merged after a reboot would sit silently in GitHub's queue
# until someone started it by hand.
#
# A *user* unit, like colima.service and platform.service, because the docker socket, the compose
# project, and runner/.env all belong to this user. Under `loginctl enable-linger` a user unit starts at
# boot with no login session — WITHOUT LINGER THIS UNIT WILL NOT RUN AT BOOT.
[Unit]
Description=platform-cicd self-hosted runner
# Requires= the docker daemon: with Colima's VM down there is no daemon for `docker compose`. After=
# platform.service because the compose file joins the `minikube` external network, which only exists
# once minikube has started (`docker compose up` fails outright if absent). Wants=, not Requires=, so a
# slow platform just makes Restart=on-failure retry rather than failing permanently on a start-order
# race.
After=colima.service platform.service network-online.target
Requires=colima.service
Wants=platform.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
# docker manages the container lifecycle (restart: always respawns the ephemeral runner between jobs);
# systemd only brings the project up at boot and down on stop.
ExecStart=/home/linuxbrew/.linuxbrew/bin/docker compose -f %h/git-workspace/claude-workspace/platform-cicd/runner/compose.yml up -d
ExecStop=/home/linuxbrew/.linuxbrew/bin/docker compose -f %h/git-workspace/claude-workspace/platform-cicd/runner/compose.yml down
# Retry until platform.service has created the minikube network (a cold boot is ~5 minutes).
Restart=on-failure
RestartSec=20

[Install]
WantedBy=default.target
